# root/etc/s6-overlay/s6-rc.d/init-pydio-cells-config/up
/etc/s6-overlay/s6-rc.d/init-pydio-cells-config/run