# root/etc/s6-overlay/s6-rc.d/svc-pydio-cells/finish
#!/usr/bin/with-contenv bash

killall -lq cells
